# quadlets/immich/immich-redis.container
[Container]
ContainerName=%N
HealthCmd=redis-cli ping || exit 1
Image=docker.io/valkey/valkey@sha256:546304417feac0874c3dd576e0952c6bb8f06bb4093ea0c9ca303c73cf458f63
Pod=immich.pod

[Service]
Restart=always

[Install]
WantedBy=default.target
