# quadlets/immich/immich-redis.container
[Container]
ContainerName=%N
DropCapability=NET_RAW
HealthCmd=redis-cli ping || exit 1
Image=docker.io/valkey/valkey@sha256:3eeb09785cd61ec8e3be35f8804c8892080f3ca21934d628abc24ee4ed1698f6
Network=immich-redis.network
NoNewPrivileges=yes
User=1000:1000
Volume=immich-redis.volume:/data:U,Z

[Service]
Restart=always
