# quadlets/immich/immich-redis.container
[Container]
ContainerName=%N
HealthCmd=redis-cli ping || exit 1
Image=docker.io/valkey/valkey@sha256:81db6d39e1bba3b3ff32bd3a1b19a6d69690f94a3954ec131277b9a26b95b3aa
Pod=immich.pod

[Service]
Restart=always

[Install]
WantedBy=default.target
