# quadlets/immich/immich-redis.container
[Container]
ContainerName=%N
HealthCmd=redis-cli ping || exit 1
Image=docker.io/valkey/valkey@sha256:a137a2b60aca1a75130022d6bb96af423fefae4eb55faf395732db3544803280
Pod=immich.pod

[Service]
Restart=always

[Install]
WantedBy=default.target
